# debian/glance-api.prerm
#!/bin/sh

set -e

case $1 in
	remove|purge)
		if [ -x /etc/init.d/glance-api ]; then
			invoke-rc.d glance-api stop || true
		fi
		;;
esac

#DEBHELPER#
